# .config/crontab
*/15 * * * * /usr/local/bin/offlineimap >/dev/null 2>&1
*/15 * * * * /usr/local/bin/vdirsyncer sync >/dev/null 2>&1
