# .config/crontab
*/15 * * * * offlineimap >/dev/null 2>&1
*/15 * * * * vdirsyncer sync >/dev/null 2>&1
